# weather/crontab/crontab
*/15 * * * * /var/www/html/Pi-Kitchen-Dashboard/bin/cron.sh
